# examples/deploy/gcp/systemd/truss-receipt-cleanup.timer
[Unit]
Description=Daily cleanup of old Truss demo receipts

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
